# overlay/etc/inittab
::sysinit:/bin/mount -a
::sysinit:/etc/init.d/rcS

::respawn:/sbin/getty -L -n -i -l /bin/autologin console 0 vt100

::shutdown:/etc/init.d/rcK
::shutdown:/bin/umount -a -r

::restart:/sbin/init
